# scripts/mac-mini/launchd/com.manacore.ssd-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.manacore.ssd-check</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mana/check-ssd.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/mana/Library/Logs/ssd-check.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mana/Library/Logs/ssd-check.log</string>
</dict>
</plist>
