# scripts/mac-mini/launchd/com.manacore.image-gen.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.manacore.image-gen</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mana/projects/manacore-monorepo/services/mana-image-gen/.venv/bin/python3</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>app.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>3025</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/mana/projects/manacore-monorepo/services/mana-image-gen</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/Users/mana/projects/manacore-monorepo/services/mana-image-gen/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/mana</string>
        <key>PORT</key>
        <string>3025</string>
        <key>FLUX_BINARY</key>
        <string>/Users/mana/flux2/flux</string>
        <key>FLUX_MODEL_DIR</key>
        <string>/Users/mana/flux2/model</string>
        <key>DEFAULT_STEPS</key>
        <string>4</string>
        <key>GENERATION_TIMEOUT</key>
        <string>300</string>
        <key>CORS_ORIGINS</key>
        <string>https://mana.how</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/manacore-image-gen.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/manacore-image-gen.error.log</string>
</dict>
</plist>
