# services/mana-geocoding/photon-self/photon-update.service
[Unit]
Description=Weekly Photon DB refresh from GraphHopper
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/photon-update.sh
# Don't fail loudly if it can't run — next week's timer tries again
SuccessExitStatus=0
