# services/mana-geocoding/photon-self/photon-update.timer
[Unit]
Description=Trigger weekly Photon DB refresh

[Timer]
# Sunday 03:30 — outside likely-usage hours
OnCalendar=Sun *-*-* 03:30:00
# Run on next boot if the system was off at the scheduled time
Persistent=true
# Spread across 30 min so we don't hammer GraphHopper at exactly :30:00
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
