# services/mana-stt/com.mana.vllm-voxtral.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mana.vllm-voxtral</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/mana/projects/mana-monorepo/services/mana-stt &amp;&amp; ./scripts/start-vllm-voxtral.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mana/projects/mana-monorepo/services/mana-stt</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>VLLM_PORT</key>
        <string>8100</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/mana/logs/vllm-voxtral.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mana/logs/vllm-voxtral.error.log</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
